# pdns/root/etc/supervisor/conf.d/pdns.conf
[program:pdns_server]
command=sh /start.sh
autorestart=true
